# packaging/deb/DEBIAN/prerm
#!/bin/sh
set -e

# Stop and disable the helper service before files are removed
if command -v systemctl >/dev/null 2>&1 && systemctl is-system-running --quiet 2>/dev/null; then
    systemctl stop    shutdowntimer-helper.service || true
    systemctl disable shutdowntimer-helper.service || true
fi

# Kill the main app if running
if command -v pkill >/dev/null 2>&1; then
    pkill -f ShutdownTimer || true
fi

exit 0
